# forge-llm-sidecar/scripts/forge-sidecar.service
# systemd unit template for the Forge LLM sidecar.
#
# Runs the sidecar on a home server so Forge can reach it from anywhere over
# Tailscale. Install as a per-user service (no root needed):
#
#   mkdir -p ~/.config/systemd/user
#   cp scripts/forge-sidecar.service ~/.config/systemd/user/
#   # edit WorkingDirectory / ExecStart below to match your install
#   systemctl --user daemon-reload
#   systemctl --user enable --now forge-sidecar
#   loginctl enable-linger "$USER"   # keep it running after logout
#
# Check status / logs:
#   systemctl --user status forge-sidecar
#   journalctl --user -u forge-sidecar -f

[Unit]
Description=Forge LLM Sidecar (deck recognition)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Absolute path to this checkout's forge-llm-sidecar directory.
WorkingDirectory=%h/forge/forge-llm-sidecar

# HOST=0.0.0.0 exposes the sidecar to the Tailscale interface. Do NOT
# port-forward 18970 on your router — Tailscale is the only ingress path.
Environment=HOST=0.0.0.0
Environment=PORT=18970
# The sidecar and llama.cpp share this host, so localhost still works.
Environment=LLM_BASE_URL=http://localhost:8080/v1

# Use the venv/interpreter that has the sidecar installed (pip install -e .).
ExecStart=%h/forge/forge-llm-sidecar/.venv/bin/uvicorn app.main:app --host ${HOST} --port ${PORT}

Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
